# systemd/dbconvd.service.in
[Unit]
Description=Dropbox File Conveyor

[Service]
ExecStart=@prefix@/bin/dbconvd
Type=simple
