# debian/postinst
#!/bin/sh
# postinst script for dropbox-conveyor

case "$1" in
    configure)
      if test ! -e /etc/dbconvd.conf ; then
	  cp /usr/share/dropbox-conveyor/dbconvd-sample.conf /etc/dbconvd.conf
      fi
      /bin/systemctl daemon-reload
      /bin/systemctl enable dbconvd.service
      /bin/systemctl start dbconvd.service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
